# template/etc/cron.hourly/rrd
#!/bin/bash
/opt/rpi-sensors/bin/net/rrd.pl
