# template/etc/cron.hourly/sqlsync
#!/bin/bash

/opt/rpi-sensors/bin/net/sqlsync.pl
